# launchd/com.crandler.claude-ha-bridge.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.crandler.claude-ha-bridge</string>

    <key>ProgramArguments</key>
    <array>
        <string>__HOME__/.config/claude-ha-bridge/venv/bin/python3</string>
        <string>__REPO__/bin/claude-ha-daemon.py</string>
    </array>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
    </dict>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <dict>
        <key>SuccessfulExit</key>
        <false/>
        <key>Crashed</key>
        <true/>
    </dict>

    <key>ThrottleInterval</key>
    <integer>10</integer>

    <key>ProcessType</key>
    <string>Background</string>

    <key>WorkingDirectory</key>
    <string>__HOME__</string>

    <key>StandardOutPath</key>
    <string>__HOME__/.config/claude-ha-bridge/daemon.log</string>

    <key>StandardErrorPath</key>
    <string>__HOME__/.config/claude-ha-bridge/daemon.log</string>
</dict>
</plist>
